# rootfs-overlay/droid4/etc/init.d/droid4-lcd
#!/sbin/openrc-run
# Copyright 1999-2018 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

depend() {
	after eudev
}

start() {
	ebegin "Setting LCD backlight to 200"
	echo 200 > /sys/class/backlight/lcd/brightness
}

stop() {
	ebegin "Nope."
}
